# systemd/publicip-collect.service
[Unit]
Description=Collect Cloudflare DNS records into DB (CLI-things publicip)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/cli-things
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/cli-things/publicip.conf
TimeoutStartSec=120
ExecStart=/usr/bin/env bash -lc './bin/publicip --collect-cf --cf-timeout 40s'

[Install]
WantedBy=multi-user.target
