# systemd/publicip-sync.timer
[Unit]
Description=Every-minute Cloudflare DNS sync (CLI-things publicip)

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=30s
Persistent=true
Unit=publicip-sync.service

[Install]
WantedBy=timers.target
